# etc/systemd/system/mnt-backup.automount
[Unit]
Description=Automount backup destination folder
After=network-online.target

[Automount]
Where=/mnt/backup
TimeoutIdleSec=480

[Install]
WantedBy=multi-user.target
